# profile/woop/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Exec=bash -c "sleep 3 && redshift -O 3500"

