# profile/woop/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Exec=redshift -O 3750

